# deploy/files/n8n-backup.timer
[Unit]
Description=n8n Hourly Backup Timer

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
