# systemd/monitoring-stack.service
[Unit]
Description=Monitoring stack (podman compose)
Requires=n150-net.service
After=n150-net.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/n8n
ExecStart=/usr/bin/podman compose -f compose/monitoring.yml up -d
ExecStop=/usr/bin/podman compose -f compose/monitoring.yml down
TimeoutStartSec=0
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
